# deploy/tg-ticket-monitor.service
[Unit]
Description=tkt.ge Telegram Ticket Monitor — checks train ticket availability and notifies via Telegram
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=tg-ticket-mon
Group=tg-ticket-mon

WorkingDirectory=/root/tg-ticket-monitor

EnvironmentFile=/root/tg-ticket-monitor/.env
ExecStart=/root/tg-ticket-monitor/.venv/bin/python3 /root/tg-ticket-monitor/bot.py

Restart=on-failure
RestartSec=10

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/root/tg-ticket-monitor/data

StandardOutput=journal
StandardError=journal

SyslogIdentifier=tg-ticket-monitor

[Install]
WantedBy=multi-user.target
